# scripts/crontab
* */1 * * * /app/.venv/bin/python3 /app/manage.py cleanup_snippets 2>&1
